# datareader.conf
[program:datareader_script]
command=/usr/local/bin/datareader
autostart=true
autorestart=true
stderr_logfile=/var/log/datareader.err.log
stdout_logfile=/var/log/datareader.out.log
